# scripts/2hjs-tracker.service
[Unit]
Description=2-Hour Job Search Tracker
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/evie-marie/Projects/2hjs-tracker
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
